# assets/supervisor/conf.d/cron.conf
[program:cron]
priority = 4
command = cron
autostart = true
user = root
redirect_stderr=true
stdout_logfile = /data/toran-proxy/logs/cron/supervisor.log
loglevel = warn
