# assets/supervisor/conf.d/php-fpm.conf
[program:php-fpm]
priority = 1
command = /usr/sbin/php-fpm7.0 -F
autostart = true
user = root
redirect_stderr=true
stdout_logfile = /data/logs/php-fpm/supervisor.log
loglevel = warn
